# deploy/kammerz.service
[Unit]
Description=Kammerz film catalog
After=network.target

[Service]
Type=simple
User=kammerz
WorkingDirectory=/opt/kammerz
EnvironmentFile=/opt/kammerz/.env
ExecStart=/opt/kammerz/kammerz
Restart=on-failure
StartLimitIntervalSec=300
StartLimitBurst=5
ProtectSystem=strict
ReadWritePaths=/opt/kammerz/data
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
