# deploy/kammerz.service
[Unit]
Description=Kammerz film catalog
After=network.target
# Restart rate-limit lives in [Unit] (systemd v229+); ignored if placed in [Service].
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=kammerz
WorkingDirectory=/opt/kammerz
EnvironmentFile=/opt/kammerz/.env
ExecStart=/opt/kammerz/kammerz
Restart=on-failure
ProtectSystem=strict
ReadWritePaths=/opt/kammerz/data
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
